# mouseemu-git/mouseemu.service
[Unit]
Description=mouseemu

[Service]
ExecStart=/usr/bin/mouseemu -nofork
KillSignal=SIGKILL

[Install]
WantedBy=multi-user.target
